# feeds/mcu/mcu/files/uci-defaults/ipq807x
[ -e /etc/config/mcu ] && exit 0

touch /etc/config/mcu

. /lib/functions.sh
. /lib/functions/mcu.sh

board=$(board_name)

case "$board" in
cig,wf196)
	mcu_add_uci_config "nrf52833_uart" "uart" "mcuboot" \
			   "zephyr-v3.3.x__hci_uart" \
			   "mcu-enable" "/dev/ttyMSM1" "115200"
	;;
edgecore,eap102)
	mcu_add_uci_config "nrf52840_usb" "usb" "mcuboot" \
			   "zephyr-v3.3.x__hci_usb" \
			   "mcu-enable"
	;;
esac

exit 0
