# feeds/wlan-ap-consumer/wlan-ap-config/files/etc/uci-defaults/90-tip-hostname
#!/bin/sh

MAC=$(cat /sys/class/net/eth0/address | tr -d ":")
uci set system.@system[-1].hostname="OpenAp-${MAC:6}"
uci commit system

exit 0
