# feeds/wlan-ap/wlan-ap-config/files/etc/uci-defaults/99-tip-network
#!/bin/sh

uci set network.wan.type=bridge
uci set network.wan6.ifname=@wan
uci set network.wan.vlan_filtering=1
uci set network.lan.vlan_filtering=1
exit 0
