# feeds/wlan-ap/wlan-ap-config/files/etc/uci-defaults/99-tip-preserve-certs
#!/bin/sh

[ -z "$(grep opensync/certs /etc/sysupgrade.conf)" ] && {
	echo /usr/opensync/certs/ >> /etc/sysupgrade.conf
}

exit 0
