# feeds/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-hostname
#!/bin.sh

hname=$(cat /sys/class/net/eth0/address | tr -d : | awk '{print tolower($0)}')
uci set system.@system[-1].hostname=$hname
uci set ucentral.config.serial=$hname
