# feeds/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-qos
#!/bin/sh

uci delete qosify.wan
uci delete qosify.wandev
uci set qosify.@defaults[-1].defaults=/tmp/qosify.conf
