# feeds/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-uhttpd
#!/bin/sh
/etc/init.d/rpcd disable
uci delete uhttpd.main
