# feeds/qca-wifi-7/ipq53xx/base-files/etc/hotplug.d/iface/85-greqos
# /etc/hotplug.d/iface/85-greqos
[ "$ACTION" == "ifup" ] || exit 0

. /lib/functions.sh

case "$(board_name)" in
cig,wf189|\
cig,wf189w|\
cig,wf189h|\
cig,wf672)
    case "$INTERFACE" in
      gre*)
        dev=$(ubus call network.interface.$INTERFACE status | jsonfilter -e '@.l3_device')
        [ -n "$dev" ] && {
          tc qdisc del dev $dev root 2>/dev/null
          tc qdisc add dev $dev root noqueue
        }
        ;;
    esac
esac