# feeds/ucentral/natlog/files/etc/init.d/natlog
#!/bin/sh /etc/rc.common
# OpenWrt init script for natlog kernel module

START=15
STOP=90

USE_PROCD=1

modfile="/lib/modules/$(uname -r)/natlog.ko"

start_service() {
    enabled=$(uci get natlog.@defaults[0].enabled 2>/dev/null)

    if [ "$enabled" = "1" ]; then
        if [ -f "$modfile" ]; then
            echo "Loading natlog kernel module..."
            insmod "$modfile" || {
                echo "Failed to load $modfile"
                return 1
            }
        else
            echo "Kernel module not found: $modfile"
            return 1
        fi
    else
        echo "natlog disabled in UCI config"
    fi
}

stop_service() {
    if lsmod | grep -q "^natlog"; then
        echo "Unloading natlog kernel module..."
        rmmod natlog
    fi
}