# feeds/ucentral/ucentral-client/files/etc/uci-defaults/zzz-ucentral
#!/bin/sh

[ -d /etc/config-shadow/ ] && exit 0

uci commit

cp -r /etc/config/ /etc/config-shadow/
rm /etc/config/network
cp /rom/etc/config/dhcp /etc/config/
config_generate
