# feeds/qca-wifi-6/ipq807x/base-files/etc/hotplug.d/net/40-no-gro
#!/bin/sh
[ add = "$ACTION" -a -n "$DEVICENAME" ] && {
	logger disable GRO on $DEVICENAME
	ethtool -K $DEVICENAME gro off
}
