# feeds/wlan-ap/opensync/files/etc/hotplug.d/net/40-rlimit-ssid
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0
rlimit=$(uci get wireless.${DEVICENAME}.rlimit)

/etc/init.d/nft-qos restart
exit 0
