# feeds/wlan-ap/opensync/files/etc/hotplug.d/net/60-dynamic-vlan
[ "$ACTION" = "add" ] || exit 0

dvid=`echo $DEVICENAME | awk -F "." '{ printf $2 }'`
[ -z "$dvid" -o "$dvid" = 0 -o "$dvid" = 1 ] && exit 0

iface=`echo $DEVICENAME | awk -F "." '{ printf $1 }'`
iswlan=`echo $iface | grep -ic wlan`
[ "$iswlan" = 0 ] && exit 0

logger -t "Add Dynamic VLAN" "dev" "$DEVICENAME" "dvid=" "$dvid"

brctl delif br-wan.$dvid $DEVICENAME
brctl addif br-wan $DEVICENAME

bridge vlan add vid $dvid dev br-wan self
bridge vlan add vid $dvid dev $(uci get network.wan.ifname)
bridge vlan add pvid $dvid vid $dvid dev $DEVICENAME untagged
bridge vlan add vid $dvid dev br-lan self

exit 0
