# feeds/ipq807x_v5.4/ipq50xx/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case "$(board_name)" in
	hfcl,ion4xi_w)
		fw_setenv boot_count 0
		;;		
	esac
}
