# feeds/ipq807x_v5.4/ipq50xx/base-files/etc/init.d/hdwdt
#!/bin/sh /etc/rc.common

START=09

boot() {
	. /lib/functions.sh

	case "$(board_name)" in
	sonicfi,rap630e)
		/sbin/hwwatchdog.sh &
		;;
	esac
}
