# feeds/mtk-app/tops-tool/files/tops-tool.init
#!/bin/sh /etc/rc.common

#the priority of TOPS driver is 51
START=52

USE_PROCD=1
NAME=tops-tool
PROG=/usr/sbin/tops-tool

start_service() {
	for script in /etc/init.d/tops-tool*; do
		if [[ -f "$script" && -x "$script" ]]; then
			source "$script"
		fi
	done
}

stop_service() {
	service_stop "${PROG}"
}
