# feeds/ucentral/ucentral-defaults/files/etc/uci-defaults/99-ucentral-opennds
#!/bin/sh

uci set opennds.@opennds[-1].enabled=0
