# feeds/ucentral/ucentral-defaults/files/etc/uci-defaults/99-ucentral-network
#!/bin/sh

uci set network.wan6.ifname=@wan
uci set network.wan.type=bridge
