# feeds/ucentral/ucentral-defaults/files/etc/uci-defaults/99-ucentral-network
#!/bin/sh

uci set network.wan6.ifname=@wan
[ -n "$(uci get network.bridge)" ] && uci set network.wan.type=bridge
