# feeds/ipq807x_v5.4/ipq60xx/base-files/etc/hotplug.d/net/macs
#!/bin/sh
mac=$(cat /etc/board.json | jsonfilter -e '@["network_device"]["'$DEVICENAME'"]'.macaddr)
[ -n "$mac" ] && ip link set $DEVICENAME address $mac
