# feeds/qca-wifi-7/cig-platform-pkg/cig-device-boot/files/cig-device.init
#!/bin/sh /etc/rc.common

START=99

boot(){
    [ -e /dev/ttyUSB0 ] && stty -F /dev/ttyUSB0 115200 cs8 -cstopb -parenb -icrnl -onlcr
}
