# feeds/ipq95xx/ipq53xx/base-files/etc/init.d/ath12k
#!/bin/sh /etc/rc.common

START=98

boot() {
	rm /etc/config/wireless
	wifi config
	cp /etc/config/wireless /etc/config-shadow
}
