# feeds/qca-wifi-7/ipq54xx/base-files/etc/uci-defaults/35-packet-steering
#!/bin/sh
# QCA IPQ platforms: default packet_steering to disabled.
# Hardware acceleration (PPE/NSS/SFE) handles packet distribution;
# software RPS provides no benefit and may cause uneven CPU load.
# Only set defaults if not already configured (preserves user settings).

changed=0

if [ -z "$(uci -q get network.globals.packet_steering)" ]; then
    uci set network.globals.packet_steering='0'
    changed=1
fi

if [ -z "$(uci -q get network.globals.steering_flows)" ]; then
    uci set network.globals.steering_flows='0'
    changed=1
fi

[ "$changed" = "1" ] && uci commit network

exit 0