# feeds/mediatek-sdk/mediatek/mt7986/base-files/etc/hotplug.d/ieee80211/11_fix_wifi_mac
[ "$ACTION" == "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
	senao,iap4300m)
		hw_mac_addr=$(mtd_get_mac_ascii u-boot-env ethaddr)
		[ "$PHYNBR" = "0" ] && macaddr_add $hw_mac_addr 1 > /sys${DEVPATH}/macaddress
		[ "$PHYNBR" = "1" ] && macaddr_add $hw_mac_addr 2 > /sys${DEVPATH}/macaddress
	;;
esac
