# feeds/tip/firstcontact/files/etc/uci-defaults/zzz-firstcontact
#!/bin/sh
[ -f "/etc/ucentral/gateway.json" ] && /etc/init.d/firstcontact disable
[ -f "/etc/ucentral/gateway.json" ] || /etc/init.d/ucentral disable
