# feeds/mtk-app/ephy-utils/files/ephy-utils.init
#!/bin/sh /etc/rc.common

START=19
STOP=20

NAME=ephy-utils

boot() {
	NEW_PATH="/usr/sbin/ephy-utils"
	PROFILE_PATH=$(grep '^export PATH=' /etc/profile | sed 's/export PATH="\(.*\)"/\1/')
	case "$PROFILE_PATH" in
		*"$NEW_PATH"*) ;;
		*) sed -i "s|export PATH=\"\(.*\)\"|export PATH=\"\1:$NEW_PATH\"|" /etc/profile
		   ;;
	esac
}
