# feeds/wlan-ap/luci/luci-theme-tip/root/etc/uci-defaults/30_luci-theme-tip
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.tip=/luci-static/tip
		set luci.main.mediaurlbase=/luci-static/tip
		commit luci
	EOF
fi

exit 0
