# feeds/ucentral/ucentral-tools/files/etc/hotplug.d/iface/90-ip-collide
[ ifup = "$ACTION" ] && {
	ubus list ucentral
	if [ $? -eq 0 ]; then
		ubus call ucentral ip_collide
	else
		/usr/share/ucentral/ip-collide.uc
	fi
}
