# feeds/ucentral/ucentral-tools/files/etc/hotplug.d/iface/90-ip-collide
[ ifup = "$ACTION" ] && {
	ip-collide
	[ $? -eq 0 ] || ubus call ucentral send '{"error": "ip/domain collision detected"}'
}
