# feeds/ucentral/ucentral-tools/files/etc/hotplug.d/iface/90-ip-collide
[ ifup = "$ACTION" ] && /usr/share/ucentral/ip-collide.uc
