# feeds/qca-wifi-7/mac80211/files/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

. /lib/functions.sh

[ "${ACTION}" = "add" ] && {
	/sbin/wifi config
}
