# feeds/qca-wifi-6/ipq60xx/base-files/etc/init.d/wdt
#!/bin/sh /etc/rc.common

START=10

boot() {
	. /lib/functions.sh

	case "$(board_name)" in
	cig,wf188|\
	cig,wf188n)
		watchdog -t 5 /dev/watchdog1
		;;
	esac
}
