# feeds/morse/netifd-morse/etc/uci-defaults/aaa-fix-phy0-to-morse

. /lib/functions.sh
. /lib/functions/uci-defaults.sh
. /lib/functions/system.sh


board=$(board_name)

case "$board" in
edgecore,eap112)

   
        uci set wireless.radio0="wifi-device"
        uci set wireless.radio0.disabled='0'
        uci set wireless.radio0.mode='ap'
        uci set wireless.radio0.type='morse'
        uci set wireless.radio0.path='platform/11009000.spi/spi_master/spi1/spi1.0'
        uci set wireless.radio0.channel='12'
        uci set wireless.radio0.hwmode='a'
        uci set wireless.radio0.band='s1g'
        uci set wireless.radio0.htmode='8'
        uci set wireless.radio0.country='US'
        uci set wireless.radio0.channels=''
        uci commit wireless.radio0
    

esac


exit 0
