# feeds/ipq807x/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case "$(board_name)" in
	edgecore,eap101|\
	edgecore,eap102)
		fw_setenv bootcount 0
		;;
	esac
}
