# feeds/ucentral/usteer2/files/etc/init.d/usteer2
#!/bin/sh /etc/rc.common

START=99
STOP=01

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/usteer.uc
	procd_set_param respawn 3600 5 0
	procd_close_instance
}
