# feeds/mediatek-sdk/mediatek-sdk/mt7986/base-files/etc/init.d/pppq-ebl.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
NAME=pppq_ebl

start_service() {
	echo 2 > /sys/kernel/debug/hnat/qos_toggle
}
