# feeds/rtkmipsel/rtk_app/files/rtk_app.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
set -x
START=38
start() {
	#delay for wifi module load
	sleep 0.5
	/usr/sbin/rtk_txcalr -w
	iw reg set US
}
stop()
{
	echo "kill rtk_txcalr\n"
}
