# feeds/ucentral/luci/luci-theme-ucentral/root/etc/uci-defaults/30_luci-theme-ucentral
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.ucentral=/luci-static/ucentral
		set luci.main.mediaurlbase=/luci-static/ucentral
		commit luci
	EOF
fi

exit 0
