# feeds/ucentral/ucentralgw/files/etc/uci-defaults/99-ucentral-firewall
#!/bin/sh

uci set firewall.wan.input=ACCEPT
uci commit firewall

exit 0
