# feeds/ucentral/ucentralgw/files/etc/uci-defaults/99-ucentral-mdns
#!/bin/sh

uci add_list umdns.@umdns[-1].network=foo
uci commit umdns

exit 0
