# feeds/ucentral/ucentralgw/files/etc/uci-defaults/99-ucentral-certs
#!/bin/sh

[ -f /etc/ucentral/certs/server-key.pem ] && exit 0
/usr/libexec/ucentral/create_certificates.sh

exit 0
