# feeds/ucentral/udhcpsnoop/files/etc/hotplug.d/iface/30-dhcpsnooping
#!/bin/sh
[ "$ACTION" == "ifup" -o "$ACTION" == "ifupdate"  ] || exit 0
[ "$(uci get dhcpsnooping.@snooping[-1].$INTERFACE)" -eq 1 ] && udhcpsnoop $INTERFACE $DEVICE
