# feeds/wlan-ap/apc/files/apc.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=40
STOP=40

PROG="/usr/bin/apc"

start_service() {
    procd_open_instance
    echo "Starting APC"
    procd_set_param command ${PROG}
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    echo "Killing APC"
    killall -s SIGTERM apc
}
