# feeds/qca/mac80211/files/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && {
	/sbin/wifi config
	[ -d /etc/config-shadow ] && cp /etc/config/wireless /etc/config-shadow/wireless
}
